# lib/configurator/roles/bash_history_backup/files/backup-bash-history.timer
[Unit]
Description=Hourly backup of bash history

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
