# lib/configurator/roles/bash_history_backup/files/backup-bash-history.service
[Unit]
Description=Backup bash history

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-bash-history
