# lib/workflows/workd/pkg/debian/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop workd || true
systemctl disable workd || true
