# lib/flow/workd/pkg/debian/etc/systemd/system/workd.service
[Unit]
Description=workd service
After=network.target

[Service]
ExecStart=/usr/bin/workd serve
Restart=on-failure

[Install]
WantedBy=multi-user.target
